# override/etc/services.d/00-registry/run
#!/usr/bin/with-contenv sh
exec s6-setuidgid guardian /usr/bin/registry serve /mnt/config/etc/registry/config.yml